# app/cron.d/crontab
0 */5 * * * /bin/bash -lc /app/cleanup_tmp.sh >> /var/log/monit.log 2>&1
